# docker/crontab
# execute the scheduler every minute
* * * * * /var/www/html/docker/requestUrl.sh GET /schedule
